# rootfs/etc/services.d/netbox-rq/run
#!/usr/bin/with-contenv execlineb 
s6-setuidgid netbox
cd /opt/netbox
/opt/netbox/venv/bin/python3
/opt/netbox/netbox/manage.py
rqworker
