# rootfs/etc/services.d/nginx/run
#!/usr/bin/execlineb -P
nginx -g "daemon off; pid /run/nginx.pid;" -c /data/nginx/nginx.conf
